# buildroot-external-khadas/board/khadas/edge2/rootfs-overlay/usr/lib/systemd/system/bluetooth-khadas.service
[Unit]
Description=Khadas Bluetooth service

[Service]
Type=forking
ExecStartPre=/bin/sh -c "/bin/echo 1 > /sys/class/rfkill/rfkill0/state"
ExecStart=/usr/bin/hciattach -s 115200 /dev/ttyS9 bcm43xx 1500000

[Install]
WantedBy=hassos-hardware.target
